# flask_util/jass_service.conf
[program:jass]
command=/home/localadmin/jass-remote/flask_util/run_flask.sh
autostart=yes
autorestart=yes
