# flask_util/jass_service.conf
[program:jass]
command=/home/localadmin/jass-remote/flask_util/run_flask.sh
directory=/home/localadmin/jass-remote/
autostart=yes
autorestart=yes
stopasgroup=true
stopsignal=QUIT
